# core/dna.bootstrap/etc/profile.d/i18n.sh
# Set up i18n variables
export LANG=es_ES.UTF-8@euro
